# src/etc/cron.d/hsmm-pi
# /etc/cron.d/hsmm-pi: CRON jobs for the HSMM-Pi project

* *  * * *     root   bash -c "[ -e /var/run/gpsd.sock ] && [ ! -e /var/run/read_gps_coordinates.lock ] && /usr/local/bin/read_gps_coordinates.pl" 2>/dev/null
